# configs/ssh/allow-ssh.service
[Unit]
Description=Allow SSH from home network
After=network.target

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'iptables -C INPUT -s 192.168.4.0/24 -p tcp --dport 22 -j ACCEPT 2>/dev/null || iptables -I INPUT 1 -s 192.168.4.0/24 -p tcp --dport 22 -j ACCEPT'
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
